# cron/crontab
@hourly command -v newsboat > /dev/null && newsboat -x reload
